# wake-up.service
[Unit]
Description=Run script on woke.

[Service]
Type=oneshot
ExecStart=/root/on-wakeup.sh

[Install]
WantedBy=suspend.target
